# sensu-server/config/supervisord.conf
[supervisord]
nodaemon=true

[program:redis]
command=/etc/init.d/redis start

[program:rabbitmq-server]
command=/etc/init.d/rabbitmq-server start

[program:sensu-server]
command=/etc/init.d/sensu-server start

[program:uchiwa]
command=/etc/init.d/uchiwa start

[program:sensu-api]
command=/etc/init.d/sensu-api start
