# tor/supervisord/privoxy-supervisor.conf
[program:privoxy]
command=/usr/sbin/privoxy --no-daemon /opt/privoxy/config
autostart = true
redirect_stderr = true
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups = 10
