# tor/supervisord/tor-supervisor.conf
[program:tor]
command=/usr/sbin/tor -f /opt/tor/torrc
autostart = true
redirect_stderr = true
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups = 10
